# crontab
0 2 * * 0 /usr/local/bin/python /app/monitor.py

